# Library/LaunchAgents/com.user.yabai-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.yabai-watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>LAUNCHD=1 /Volumes/m2/Usuarios/.local/bin/yabai-check.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>120</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/yabai-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/yabai-watchdog.err</string>
</dict>
</plist>
